# deploy/pagescms/pagescms-backup.service
[Unit]
Description=Back up the Pages CMS PostgreSQL database
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/local/sbin/pagescms-backup-database
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/pagescms
RestrictAddressFamilies=AF_UNIX
